# templates/default/livy.service.erb
[Unit]
Description=Livy server
After=syslog.target network.target remote-fs.target <%= @deps %>

[Service]
User = <%= node['livy']['user'] %>

PIDFile=<%= node['livy']['pid_file'] %>
ExecStartPre=/bin/bash -c 'while ! /usr/bin/host <%= @rm_rpc_endpoint %> ; do /bin/sleep 5; done'
ExecStart=<%= node['livy']['base_dir'] %>/bin/start-livy.sh
ExecStop=<%= node['livy']['base_dir'] %>/bin/stop-livy.sh
LimitNOFILE=20000
TimeoutSec= 10min

Type=forking

[Install]
WantedBy=multi-user.target
